# deploy/systemd/kai-service-watchdog.service
[Unit]
Description=KAI Pi External Service Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/kai/ai_analyst_trading_bot
EnvironmentFile=/home/kai/ai_analyst_trading_bot/.env
Environment=PYTHONIOENCODING=utf-8
ExecStart=/bin/bash /home/kai/ai_analyst_trading_bot/scripts/pi_service_watchdog.sh
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/kai/ai_analyst_trading_bot
