# deploy/systemd/kai-auto-annotate-blocked.timer
[Unit]
Description=Trigger KAI Auto-Annotate Blocked every 6 hours (D-227)
Requires=kai-server.service

[Timer]
# 45 min after boot first-run, then every 6h. Offset from kai-auto-annotate.timer
# so CoinGecko calls do not stack immediately after boot.
OnBootSec=45min
OnUnitActiveSec=6h
AccuracySec=2min
Unit=kai-auto-annotate-blocked.service

[Install]
WantedBy=timers.target
