# deploy/systemd/kai-recalc-cycle.service
[Unit]
Description=KAI Recalc-Cycle — daily refresh of source-reliability/lifecycle/bayes/confluence/ph5
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ai_analyst_trading_bot
Environment=PYTHONPATH=/home/ubuntu/ai_analyst_trading_bot
EnvironmentFile=-/home/ubuntu/ai_analyst_trading_bot/.env

# Sequenz: jeder Script-Lauf ist eigenständig; "-" prefix sorgt dafür dass
# ein Failure den nächsten Recalc nicht killt. Operator sieht Failures im
# Journal und über die kai-pi-health.timer-Output-mtime-Alarme.
ExecStart=-/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/source_reliability_recalc.py
ExecStart=-/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/source_lifecycle_recalc.py
ExecStart=-/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/source_confluence_recalc.py
ExecStart=-/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/bayes_posterior_recalc.py
ExecStart=-/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/ph5_feature_analysis_recalc.py

StandardOutput=journal
StandardError=journal
TimeoutStartSec=900

# Resource caps — recalc-cycle is read-heavy (alert_audit.jsonl, DB joins).
# Each script holds the full alert_audit in memory; 512M cap reflects
# observed peak at ~80k audit rows.
MemoryMax=512M
CPUQuota=80%
