# deploy/systemd/kai-paper-trading.timer
[Unit]
Description=Trigger KAI Paper-Trading every 10 min
# NO `Requires=kai-server.service` here ON PURPOSE. `Requires=` on a *timer* makes
# `systemctl stop kai-server` deactivate this timer too (Requires= stop-propagation)
# and it is NOT restarted when the server returns — silently dead until reboot. The
# kai-server dependency lives on the paired .service (After= + Requires=), which still
# gates the actual work. Do not re-add Requires= here.

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Unit=kai-paper-trading.service

[Install]
WantedBy=timers.target
