# deploy/systemd/kai-real-analysis-paper-feed.service
# KAI Real-Analysis Paper Feeder (Goal 2026-06-10).
#
# Fail-closed: the underlying feeder is a no-op unless the three-arm override is
# fully armed (REAL_ANALYSIS_PAPER_ENABLED + ..._ALLOW_PAPER_WHILE_ENTRY_DISABLED
# + ..._ENTRY_DISABLED_OVERRIDE_ACK == sentinel). Installed but NOT enabled by
# default (see pi_install_systemd.sh: present in UNITS, absent from
# ENABLE_ON_INSTALL). Enabling the .timer is a deliberate operator step.
[Unit]
Description=KAI Real-Analysis Paper Feeder — inject real eligible analyses as paper trades (fail-closed, disarmed by default)
After=network-online.target kai-server.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ai_analyst_trading_bot
Environment=PYTHONPATH=/home/ubuntu/ai_analyst_trading_bot
EnvironmentFile=-/home/ubuntu/ai_analyst_trading_bot/.env

ExecStart=/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/real_analysis_paper_feed.py

StandardOutput=journal
StandardError=journal
