# deploy/systemd/kai-server.service
[Unit]
Description=KAI Pipeline Server (FastAPI + in-process schedulers)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/kai/ai_analyst_trading_bot
EnvironmentFile=/home/kai/ai_analyst_trading_bot/.env
# D-208: systemd opens StandardOutput=append: as root before dropping to User=,
# creating root-owned log files that block the service-user from rotating/reading.
# `+` prefix runs the chown as root regardless of User=.
ExecStartPre=+/bin/bash -c 'mkdir -p /home/kai/ai_analyst_trading_bot/logs && chown -R ubuntu:ubuntu /home/kai/ai_analyst_trading_bot/logs'
# --timeout-graceful-shutdown 15: uvicorn-Default wartet UNENDLICH auf das
# Schließen offener Verbindungen. Hinter cloudflared (Keep-Alive zur Origin)
# beendet sich uvicorn dadurch nie → systemd SIGKILL erst nach TimeoutStopSec
# (= 90s Default) → jeder Restart riss ein ~90s-Offline-Loch (CF-520 im
# Dashboard). Mit 15s schließt uvicorn hängende Verbindungen sauber und exitet
# (verifiziert 2026-06-23: 15s "Deactivated successfully" statt 90s SIGKILL).
ExecStart=/home/kai/ai_analyst_trading_bot/.venv/bin/python -m uvicorn app.api.main:app --host 127.0.0.1 --port 8000 --log-level info --timeout-graceful-shutdown 15
Restart=on-failure
RestartSec=10
# Backstop: kappt den Stop hart bei 20s, falls uvicorn doch mal nicht sauber
# beendet (App-State ist append-only/JSONL → SIGKILL ist seit jeher survivable).
TimeoutStopSec=20s
StandardOutput=append:/home/kai/ai_analyst_trading_bot/logs/server.log
StandardError=append:/home/kai/ai_analyst_trading_bot/logs/server.err.log
# Hardening (RPi OS Bookworm supports all of these)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/kai/ai_analyst_trading_bot

[Install]
WantedBy=multi-user.target
