# deploy/systemd/kai-health-check.service
[Unit]
Description=KAI System Health Check — Pi-side authoritative probe with Telegram-on-issue
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ai_analyst_trading_bot
Environment=PYTHONPATH=/home/ubuntu/ai_analyst_trading_bot
EnvironmentFile=-/home/ubuntu/ai_analyst_trading_bot/.env
ExecStart=/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python -m app.cli.main alerts health-check --telegram-on-issue --lookback-hours 24 --notify-cooldown-minutes 30 --exit-on-stale
StandardOutput=journal
StandardError=journal
TimeoutStartSec=90

# Resource caps — probe is read-only and short.
# MemoryMax 256M→512M (2026-06-23): der App-Import wuchs durch #387/#388 über
# 256M → der health-check lief OOM (Result=oom-kill, exit 2). Auf dem Pi war
# /etc bereits live auf 512M gehoben+verifiziert ("All systems healthy"); hier
# nachgezogen, damit ein Repo→/etc-Redeploy nicht auf 256M zurückregressiert.
MemoryMax=512M
CPUQuota=50%
