# deploy/systemd/kai-recalc-cycle.service
[Unit]
Description=KAI Recalc-Cycle — daily refresh of source-reliability/lifecycle/bayes/confluence/ph5
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ai_analyst_trading_bot
Environment=PYTHONPATH=/home/ubuntu/ai_analyst_trading_bot
EnvironmentFile=-/home/ubuntu/ai_analyst_trading_bot/.env

# Wrapper runs all 5 recalc steps best-effort (a crashed step does NOT abort the
# rest — that independence is why we used ExecStart=- per line). It then exits
# non-zero if ANY step failed, so the unit is marked `failed` immediately instead
# of reporting success while an output silently goes stale (2026-06-24: a bayes
# circular-import crash stayed invisible at the unit level for ~31h under the old
# `-` lines). Per-step stdout/stderr still flows to the journal.
ExecStart=/home/ubuntu/ai_analyst_trading_bot/.venv/bin/python /home/ubuntu/ai_analyst_trading_bot/scripts/recalc_cycle.py

StandardOutput=journal
StandardError=journal
TimeoutStartSec=900

# Resource caps — recalc-cycle is read-heavy (alert_audit.jsonl, DB joins).
# Each script holds the full alert_audit in memory; 512M cap reflects
# observed peak at ~80k audit rows.
MemoryMax=512M
CPUQuota=80%
